# sigil-hardware/services/audio-manager.service
[Unit]
Description=Sigil Audio Manager — orchestrates radio-fetcher and audio-player
After=network-online.target sigil-pulseaudio.service bt-connect.service
Wants=network-online.target sigil-pulseaudio.service bt-connect.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=sigil
Group=sigil
ExecStart=/usr/local/bin/audio-manager.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
