# sigil-hardware/services/bt-connect.service
[Unit]
Description=Sigil Bluetooth Auto Connect
# bt-connect DEBE esperar a que PulseAudio esté completamente listo para
# registrar los endpoints A2DP en BlueZ.
# Sin esto: "br-connection-profile-unavailable" / "Protocol not available"
After=bluetooth.service bluetooth.target sigil-pulseaudio.service sound.target
Wants=bluetooth.service sigil-pulseaudio.service
# Reiniciar indefinidamente si falla (la Pi puede no tener Bluetooth listo de inmediato)
StartLimitIntervalSec=0

[Service]
Type=simple
User=sigil
Group=sigil
# PulseAudio is started by sigil-pulseaudio.service; the script then waits for
# the actual A2DP endpoints instead of using a blind boot delay.
ExecStart=/usr/local/bin/bt-connect.sh
Restart=on-failure
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
