# sigil-hardware/services/radio-stream.service
[Unit]
Description=Sigil Radio Stream 24/7
# Legacy rollback only. Routing is delegated to sigil-audio-route.sh by the
# script; this unit does not preselect a sink.
After=network-online.target pulseaudio.service bt-connect.service sound.target
Wants=network-online.target pulseaudio.service bt-connect.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=sigil
Group=sigil
# PULSE_SERVER se resuelve dinámicamente en radio-stream.sh usando `id -u sigil`
# NO hardcodear el UID aquí — el script lo calcula en runtime.
ExecStartPre=/bin/sleep 20
ExecStart=/usr/local/bin/radio-stream.sh
KillMode=control-group
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
