# sigil-hardware/services/ssh-monitor.service
[Unit]
Description=Sigil SSH Monitor — security maintenance mode when SSH is active
After=network.target ssh.service
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/ssh-monitor.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
